# deploy/spaitra-ocr.service
[Unit]
Description=Spaitra OCR service
After=network.target

[Service]
Type=simple
User=spaitra
Group=spaitra
WorkingDirectory=__REPO_ROOT__
EnvironmentFile=/opt/spaitra/.ocr.env
ExecStart=__REPO_ROOT__/venv-ocr/bin/python -m services.ocr.run
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
