# deploy/spaitra-core.service
[Unit]
Description=Spaitra core backend
After=network.target

[Service]
Type=simple
User=spaitra
Group=spaitra
UMask=0077
WorkingDirectory=__REPO_ROOT__
EnvironmentFile=/opt/spaitra/.env
ExecStart=__REPO_ROOT__/venv-core/bin/gunicorn \
    -w 1 \
    -b 127.0.0.1:5000 \
    --timeout 120 \
    --access-logfile - \
    services.core.wsgi:application
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
